# nginx/supervisor.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi -s /tmp/uwsgi.sock -w wsgiapp:application

[program:nginx]
command = /usr/sbin/nginx